# selfstar-deploy/systemd/selfstar-cert-renew.service
[Unit]
Description=Renew Let's Encrypt certificates for selfstar.duckdns.org (Docker Compose)
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/selfstar/selfstar-deploy
ExecStart=/usr/bin/docker compose run --rm certbot renew --webroot -w /var/www/certbot
ExecStartPost=/usr/bin/docker compose exec nginx nginx -s reload
TimeoutStartSec=0
